# docker/quadlets/nginx-vllm.container.template
[Unit]
Description=Proxy Nginx pour vLLM - Environnement ${ENV_NAME}
After=vllm-qwen-${ENV_NAME}.service
Wants=vllm-qwen-${ENV_NAME}.service

[Container]
Image=${NGINX_IMAGE}
ContainerName=nginx-vllm-qwen-${ENV_NAME}
Network=vllm-qwen-${ENV_NAME}.network

# Port exposé sur l'hôte
PublishPort=${NGINX_HTTPS_HOST_PORT}:443

# Variables d'environnement
Environment=TZ=Pacific/Noumea
Environment=NGINX_BACKEND_UPSTREAM=vllm_qwen-${ENV_NAME}:8000
Environment=NGINX_HTTPS_LISTEN_PORT=443
Environment=NGINX_BACKEND_WARMUP_ENABLED=true
Environment=NGINX_BACKEND_WARMUP_PROMPT=warmup
Environment=NGINX_BACKEND_WARMUP_MAX_TOKENS=1
Environment=NGINX_BACKEND_WARMUP_TIMEOUT=120

# Volumes / Secrets
Volume=${RUNTIME_DIR}/secrets/NGINX_SSL_CERT:/run/secrets/NGINX_SSL_CERT:ro
Volume=${RUNTIME_DIR}/secrets/NGINX_SSL_KEY:/run/secrets/NGINX_SSL_KEY:ro
Volume=${RUNTIME_DIR}/secrets/NGINX_SSL_FULLCHAIN:/run/secrets/NGINX_SSL_FULLCHAIN:ro
Volume=${RUNTIME_DIR}/vllm_api_key.txt:/run/secrets/vllm_api_key:ro

[Service]
# Wait longer than vLLM HealthStartPeriod for large model loads.
# NOTE: use $$ to avoid envsubst eating shell variables.
TimeoutStartSec=1h
ExecStartPre=/bin/sh -c 'for i in $$(seq 1 3600); do status=$$(podman inspect -f "{{.State.Health.Status}}" vllm_qwen-${ENV_NAME} 2>/dev/null || true); [ "$$status" = "healthy" ] && exit 0; sleep 2; done; exit 1'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
